# systemd-timers/novarewind-daily.timer
[Unit]
Description=Execute task every day

[Timer]
OnCalendar=daily
Persistent=true

Unit=novarewind-daily.service

[Install]
WantedBy=timers.target
